# services/cli/systemd/longrunning.service.in
[Unit]
Description={{DESCRIPTION}}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_ROOT}}
EnvironmentFile={{ENV_FILE}}
Environment=PYTHONUNBUFFERED=1
ExecStart={{PYTHON}} -m {{SERVICE_MODULE}} {{ARGS}}

Restart=on-failure
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
KillMode=control-group
TimeoutStartSec=infinity

[Install]
WantedBy={{WANTED_BY}}