# services/systemd/longrunning.service.in
[Unit]
Description={{DESCRIPTION}}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{PROJECT_ROOT}}
ExecStart={{PYTHON}} -m {{SERVICE_MODULE}}
EnvironmentFile={{ENV_FILE}}

Restart=on-failure
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
KillMode=control-group

TimeoutStartSec=infinity

[Install]
WantedBy=default.target
