# services/podcast_transcriber/podcast_transcriber.service.in
[Unit]
Description=podcast_transcriber podcast transcribing service

[Service]
Type=oneshot
WorkingDirectory={{PROJECT_ROOT}}
ExecStart={{PYTHON}} -m {{SERVICE_MODULE}}
EnvironmentFile={{ENV_FILE}}

[Install]
WantedBy=default.target

