# services/podcast_transcriber/podcast_transcriber.timer.in
[Unit]
Description=podcast_transcriber Periodic podcast transciber

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
