# services/label_en/label_en.service.in
[Unit]
Description=label_en language labeling service

[Service]
Type=oneshot
WorkingDirectory={{PROJECT_ROOT}}
ExecStart={{PYTHON}} -m {{SERVICE_MODULE}}
EnvironmentFile={{ENV_FILE}}

[Install]
WantedBy=default.target

