# services/youtube_transcriber/youtube_transcriber.service.in
[Unit]
Description=youtube_transcriber Youtube video transcribing service

[Service]
Type=oneshot
WorkingDirectory={{PROJECT_ROOT}}
ExecStart={{PYTHON}} -m {{SERVICE_MODULE}}
EnvironmentFile={{ENV_FILE}}

[Install]
WantedBy=default.target

